# docker/s6/proxy/run
#!/command/execlineb -P
caddy run --adapter caddyfile --config /app/Caddyfile
