# docker/s6/gunicorn/run
#!/command/execlineb -P
cd /app/sortiment
gunicorn sortiment.wsgi --bind 127.0.0.1:8001 --workers 4 --max-requests 1000
