# docker/s6/migrations/up
#!/command/with-contenv sh
cd /app/sortiment
python manage.py migrate
