# docker/s6/migrations/up
/command/with-contenv sh -c "cd /app/sortiment; python manage.py migrate"
