# docker/s6/migrations/up
#!/command/execlineb -P
cd /app/sortiment
python manage.py migrate
